# src/custom/etc/rc.local
#!/bin/sh -e

# Install local debian packages
sh /root/scripts/install-packages.sh

# Start system services
sh /root/scripts/start-services.sh

exit 0